# stream/skylight-stream@.service
# Instanced TikTok/Twitch/etc pusher — one unit per stream/.env.<name>:
#   sudo cp stream/skylight-stream@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl start skylight-stream@tiktok skylight-stream@twitch
# NO enable: streams are started on purpose, not by power cuts.
[Unit]
Description=Skylight live pusher (%i)
After=network-online.target skylight-server.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/skylight/stream
ExecStart=/home/pi/skylight/stream/start-stream.sh %i
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
