# pi-setup/skylight-server.service
# systemd unit for the Skylight server (data + WebSocket + serves the built web app).
# Installed to /etc/systemd/system/ by install-on-pi.sh.
[Unit]
Description=Skylight server
After=network-online.target dump1090-fa.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APPDIR__
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=DATA_SOURCE=radio
Environment=AIRCRAFT_JSON_URL=http://localhost:8080/aircraft.json
ExecStart=__PNPM__ -F server start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
