# pi-setup/skylight-tracker.service
# systemd unit for the Skylight camera tracker (PTZ pointing + vision +
# MJPEG video + debug/TV endpoints on :3001).
# Installed to /etc/systemd/system/ by install-on-pi.sh / deploy-to-pi.sh.
[Unit]
Description=Skylight camera tracker
After=network-online.target skylight-server.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APPDIR__
Environment=NODE_ENV=production
Environment=TRACKER_PORT=3001
Environment=TRACKER_HOST=0.0.0.0
Environment=SERVER_WS=ws://localhost:3000/ws
ExecStart=__PNPM__ -F tracker start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
